# configs/openrc-unbound
#!/sbin/openrc-run
name="unbound"
description="Unbound DNS resolver"

config="/etc/unbound/unbound.conf"
pidfile="/etc/unbound/unbound.pid"

command="/usr/sbin/unbound"
command_args="-c $config"

depend() {
        need net
        use logger
        provide dns
        after auth-dns entropy
}

checkconfig() {
        ebegin "Checking $config"
        unbound-checkconf -f "$config" >/dev/null
        eend $?
}

start_pre() {
        checkconfig

}

start() {
    ebegin "Starting $name"
    start-stop-daemon --start --quiet --pidfile $pidfile --exec $command -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping $name"
    start-stop-daemon --stop --quiet --pidfile $pidfile
    eend $?
}